# systemd/mochi-frontend.service
[Unit]
Description=Mochi Next.js frontend
After=network.target

[Service]
WorkingDirectory=/root/mochi/frontend
ExecStart=/usr/bin/npm run start -- --hostname 0.0.0.0 --port 3000
Restart=always
User=root
Environment=NODE_ENV=production
EnvironmentFile=/root/mochi/frontend/.env.local

[Install]
WantedBy=multi-user.target
